# playerctl-pause.path
[Unit]
Description=Monitor sound device changes

[Path]
PathChanged=/dev/snd

[Install]
WantedBy=default.target