# playerctl-pause.service
[Unit]
Description=Pause playerctl when sound device changes
After=sound.target

[Service]
Type=oneshot
ExecStart=/usr/bin/playerctl pause